# zookeeper/config/supervisor.d/zookeeper.ini
[program:zookeeper]
command=/opt/zookeeper/bin/zkServer.sh start-foreground
redirect_stderr=true
